# scripts/lelamp.service.example
[Unit]
Description=LeLamp Runtime
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/lelamp_runtime
EnvironmentFile=/home/pi/lelamp_runtime/.env
ExecStart=/usr/bin/env bash -lc 'cd /home/pi/lelamp_runtime && exec /usr/bin/script -qefc "/home/pi/.local/bin/uv run smooth_animation.py console" /dev/null'
Restart=always
RestartSec=5
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
